# crontab
2 2 * * * /root/backup.sh > /proc/$(cat /var/run/crond.pid)/fd/1 2>&1
