# crontab
1 * * * * /root/pg_backup.sh
1 * * * * /root/fs_backup.sh
